# infra/latitude/boring-net.service
[Unit]
Description=boring computers guest networking (bridge, dnsmasq, NAT, egress firewall)
After=network-online.target dnsmasq.service
Wants=network-online.target
Before=boringd.service

[Service]
Type=oneshot
ExecStart=/opt/boring/bin/net-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
